# platform/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=debug
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:nginx]
command=/usr/local/nginx/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
autorestart=true
priority=10
redirect_stderr=true

[program:flask]
command=gunicorn -w 4 -b 0.0.0.0:5000 app:app
directory=/  
stdout_logfile=/var/log/flask-gunicorn.log
autorestart=true
priority=20
environment=PYTHONUNBUFFERED=1
redirect_stderr=true

